# scripts/com.znonymity.piasnews.social.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.znonymity.piasnews.social</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bytedance/Library/Application Support/piasnews/update_social_agent_reach.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PIASNEWS_SOCIAL_MIN_INTERVAL_SECONDS</key>
    <string>10800</string>
  </dict>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/piasnews-social.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/piasnews-social.err</string>
</dict>
</plist>
